# deploy/openarcserver.service
[Unit]
Description=OpenArcServer DX Cluster Server
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openarcserver
Group=openarcserver
WorkingDirectory=/opt/openarcserver
ExecStart=/opt/openarcserver/OpenArcServer.Server

# Always restart — .NET host exits with code 0 even on fatal background
# service exceptions, so on-failure would not trigger a restart.
Restart=always
RestartSec=10s

# Allow binding to ports below 1024 (e.g. telnet on port 23)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Log to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openarcserver

KillSignal=SIGTERM
TimeoutStopSec=30

# Hardening (keep writable paths explicit)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/opt/openarcserver/data /opt/openarcserver/logs

# Optionally override settings without editing appsettings.json:
# Environment=Server__NodeCallsign=W1AW-2
# Environment=Api__AdminKey=changeme
# Environment=Rbn__Callsign=W1AW

[Install]
WantedBy=multi-user.target
